# launchd/com.userflow.update-segments.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.userflow.update-segments</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/jack.lee/Library/Application Support/UserFlowSegments/update-segments-local.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/jack.lee/Library/Application Support/UserFlowSegments</string>

  <key>RunAtLoad</key>
  <true/>

  <key>StartInterval</key>
  <integer>1800</integer>

  <key>StandardOutPath</key>
  <string>/Users/jack.lee/Library/Logs/userflow-segments-update.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/jack.lee/Library/Logs/userflow-segments-update.err.log</string>
</dict>
</plist>
